# clipboard-sync.service
[Unit]
Description=Bidirectional clipboard sync (PRIMARY <-> CLIPBOARD)

[Service]
ExecStart=/usr/bin/clipboard-sync
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
